# crontab
# Запуск через curl на сервис app
0 * * * * curl -X GET http://app:${APP_PORT}/api/cron/run -s >> /var/log/cron.log 2>&1

# Дата для проверки
* * * * * date >> /var/log/cron.log 2>&1
